# contrib/cron.hourly
#!/bin/sh

NAME=encore
USER=encore
USERHOME=$(eval echo ~encore)
EXE=$USERHOME/go/bin/encore-parser
LOGHOME=/var/log/encore
start-stop-daemon --start --quiet \
	--pidfile /var/run/$NAME-parser.pid --make-pidfile \
	--user $USER --chuid $USER \
	--exec $EXE -- \
		-database="dbname=encore host=ns2.cs.princeton.edu sslmode=disable" \
		-geoip_database=$USERHOME/GeoIP.dat \
		-logfile=$LOGHOME/$NAME-parse.log
